# archives/rootfs/etc/profile.d/alsa_env.sh
# card 0: mini jack
# card 1: rockchiphdmi
export ALSA_PCM_CARD=1
export ALSA_PCM_DEVICE=0
