# archives/rootfs/etc/profile.d/sdl2_env.sh
# audio: alsa
# video: kmsdrm + opengles2
export SDL_AUDIODRIVER=alsa
export SDL_VIDEODRIVER=KMSDRM
export SDL_RENDER_DRIVER=opengles2
export SDL_VIDEO_GL_DRIVER=libGLESv2.so
export SDL_VIDEO_EGL_DRIVER=libEGL.so
