# archives/rootfs/etc/rc.local
#!/bin/sh -e

# init wifi
rk_wifi_init /dev/ttyS0 > /tmp/rk_wifi_init.log

# disable DP/eDP
echo off > /sys/class/drm/card0-DP-1/status
echo off > /sys/class/drm/card0-eDP-1/status

# run once script
if [ -f /etc/init.d/run_once_on_boot.sh ];then
    /etc/init.d/run_once_on_boot.sh
fi

exit 0
